# deploy/systemd/bms-backup.timer
[Unit]
Description=Run BMS backup daily at midnight

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
